# config/supervisord-flowerd.conf
; ======================================
;  flower for celery supervisor example
; ======================================
; Place under /etc/supervisor/conf.d/
[program:flowerd]
; Set full path to celery program if using virtualenv
;command=/home/vitaly/django-primes-interview/config/flowerd.sh
command=/home/vitaly/django-primes-interview/config/runinenv.sh /home/vitaly/django-primes-interview/.venv/bin/celery flower -A primes --port=5555

; Load env variables
; environment=CELERY_BASIC_AUTH=%(CELERY_BASIC_AUTH)s

; Alternatively,
;command=celery --app=your_app.celery:app worker --loglevel=INFO -n worker.%%h
; Or run a script
;command=celery.sh

directory=/home/vitaly/django-primes-interview/primes/
user=celery
numprocs=1
stdout_logfile=/var/log/celery/flower.log
stderr_logfile=/var/log/celery/flower.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000
